# remount-root-rw
#!/sbin/openrc-run
description="Remount root filesystem in read-write mode"

depend() {
    before cloud-init
}

start() {
    ebegin "Remounting root filesystem in read-write mode"
    mount -o remount,rw /
    eend $?
}
